# recipes-core/ssh-pubkey-server/files/ssh-pubkey-server-init
#!/bin/sh
### BEGIN INIT INFO
# Provides:          ssh-pubkey-server
# Required-Start:    $network $remote_fs $syslog $searcher-container
# Required-Stop:     $network $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start SSH Public Key Server daemon at boot time
# Description:       Enable SSH Public Key Server service provided by daemon.
### END INIT INFO

DAEMON=/usr/bin/httpserver
NAME=ssh-pubkey-server
DESC="SSH Public Key Server"
SSH_PUBKEY_SERVER_PORT=5001
KEY_FILE="/etc/searcher/ssh_hostkey/host_key.pub"
DAEMON_ARGS="--listen-addr=127.0.0.1:${SSH_PUBKEY_SERVER_PORT} --ssh-pubkey-file ${KEY_FILE}"
PIDFILE=/var/run/$NAME.pid
LOGFILE=/tmp/httpserver.log

PROXY_BIN=/usr/bin/proxy-server
PROXY_NAME=cvm-reverse-proxy-server
PROXY_ARGS="--listen-addr=0.0.0.0:8745 \
            --target-addr=http://localhost:${SSH_PUBKEY_SERVER_PORT} \
            --server-attestation-type=azure-tdx"
PROXY_PIDFILE=/var/run/${PROXY_NAME}.pid
PROXY_LOGFILE=/tmp/${PROXY_NAME}.log

wait_for_ssh_key() {
    echo "Waiting for SSH host key generation..."
    max_attempts=10
    attempt=1
    
    while [ $attempt -le $max_attempts ]; do
        echo "Attempt $attempt/$max_attempts: Checking for SSH host key..."
        
        # Check if file exists
        if [ -f "$KEY_FILE" ]; then
            # Verify it's a valid SSH public key by checking for "ssh-ed25519" prefix
            if grep -q "^ssh-ed25519 " "$KEY_FILE"; then
                echo "Valid SSH host key found"
                return 0
            else
                echo "File exists but doesn't appear to be a valid SSH key"
            fi
        else
            echo "Host key file not found yet"
        fi
        
        sleep 10
        attempt=$((attempt + 1))
    done
    
    echo "SSH host key was not generated within timeout"
    return 1
}

start_proxy() {
    echo -n "Starting CVM Reverse Proxy: "
    start-stop-daemon -S \
    --make-pidfile \
    --pidfile "$PROXY_PIDFILE" \
    -m -b \
    --exec /bin/sh -- -c "exec $PROXY_BIN $PROXY_ARGS >> $PROXY_LOGFILE 2>&1"
    echo "Started CVM Reverse Proxy."
}

start() {
  # Wait for SSH key to be generated by the container
  if ! wait_for_ssh_key; then
      echo "Failed to find valid SSH host key, cannot start server"
      return 1
  fi
  
  # Start the CVM reverse proxy server
  start_proxy
  
  # Start the SSH public key server
  echo "Starting $DESC: "
  start-stop-daemon -S \
        --make-pidfile \
        --pidfile "$PIDFILE" \
        -m -b \
        --exec /bin/sh -- -c "exec $DAEMON $DAEMON_ARGS 2>&1 | tee -a $LOGFILE"
    echo "$NAME started."
}

stop() {
  echo "Stopping $DESC..."
    start-stop-daemon --stop --quiet --pidfile "$PIDFILE" || true
    rm -f "$PIDFILE"

    echo "Stopping CVM Reverse Proxy Server..."
    start-stop-daemon --stop --quiet --pidfile "$PROXY_PIDFILE" || true
    rm -f "$PROXY_PIDFILE"

    echo "$NAME and reverse proxy stopped."
}

restart() {
  echo "Restarting $DESC: "
  stop
  sleep 1
  start
}

test -x $DAEMON || exit 0

set -e

case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  restart|force-reload)
    restart
    ;;
  *)
    N=/etc/init.d/$NAME
    echo "Usage: $N {start|stop|restart|force-reload}" >&2
    exit 1
    ;;
esac

exit 0
