# recipes-core/searcher-ssh-key/files/searcher-ssh-key.sh.in
#!/bin/sh
### BEGIN INIT INFO
# Provides:          searcher-ssh-key
# Required-Start:    $all
# Required-Stop:     
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6   
# Short-Description: Creates searcher user and adds searcher ssh key
### END INIT INFO

start() {
  # Start the attested TLS reverse proxy server
  start_proxy

  echo "Starting searcher-ssh-key script"
  # Add searcher user
  if ! id "searcher" &>/dev/null; then
      adduser --disabled-password --gecos "" searcher
      echo "Created searcher user"
  else
      echo "searcher user already exists"
  fi

  # Create .ssh directory and set permission
  mkdir -p /home/searcher/.ssh
  chmod 700 /home/searcher/.ssh
  chown searcher:searcher /home/searcher/.ssh

  # Add searcher user to subuid & subgid
  echo "searcher:100000:65536" > /etc/subuid
  echo "searcher:100000:65536" > /etc/subgid

  # Add searcher's ssh key to authorized_keys file
  SSH_KEY="@SEARCHER_SSH_KEY@"
  echo "$SSH_KEY" > /home/searcher/.ssh/authorized_keys
  chmod 600 /home/searcher/.ssh/authorized_keys
  chown searcher:searcher /home/searcher/.ssh/authorized_keys
  echo "Added Searcher SSH key to searcher user's authorized_keys"

  # Generate dropbear host public key to share upon attestation
  dropbearkey -y -f /etc/dropbear/dropbear_rsa_host_key | grep -m 1 "^ssh-rsa " > /etc/dropbear/dropbear_rsa_host_key.pub
}

case "$1" in
  start)
    start
    ;;
  stop)
    echo "Stopping searcher-ssh-key script"
    ;;
  *)
    echo "Usage: $0 {start|stop}"
    exit 1
    ;;
esac

exit 0
